# ansible/roles/webapp/templates/gunicorn.service.j2
[Unit]
Description={{app_name}} HTTP gunicorn

[Service]
User={{user_name or app_name}}
KillSignal=SIGQUIT
ExecStart={{src_root}}/manage.py run_gunicorn

[Install]
WantedBy=multi-user.target
